# roles/web/files/borgback.timer
[Unit]
Description=Timer For BorgBack service

[Timer]
OnUnitActiveSec=15m

[Install]
WantedBy=multi-user.target
