# roles/monitor/files/systemd-journal-remote.timer
[Unit]
Description=Timer For system-journal-remote

[Timer]
OnActiveSec=1s

[Install]
WantedBy=multi-user.target
