# roles/mysql/files/fullbackup.timer
[Unit]
Description=Timer For MySQL Full Backup service

[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=multi-user.target
