# roles/mysql/files/incbackup.service
[Unit]
Description=MySQL Incremental Backup Script

[Service]
ExecStart=/bin/bash /root/inc-backup.sh

[Install]
WantedBy=multi-user.target
