# roles/mysql/files/incbackup.timer
[Unit]
Description=Timer For MySQL Incremental Backup service

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=multi-user.target
